# Library/LaunchDaemons/com.mount.jobs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.mount.jobs</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/desn-render-jobs.sh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
<!--	<key>StartInterval</key>
	<integer>10</integer>
-->	<key>StandardOutPath</key>
	<string>/var/log/com.deadline/desn-render-jobs-daemon.log</string>
	<key>StandardErrorPath</key>
	<string>/var/log/com.deadline/desn-render-jobs-daemon_err.log</string>
</dict>
</plist>
