# Library/LaunchAgents/com.mount.dr.user.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.mount.dr.user</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/desn-render-dr-user.sh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/var/log/com.deadline/desn-render-dr-user.log</string>
	<key>StandardErrorPath</key>
	<string>/var/log/com.deadline/desn-render-dr-user_err.log</string>
</dict>
</plist>
